# pi/systemd/plant-eye.service
[Unit]
Description=Plant Eye Camera Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
EnvironmentFile=/etc/plant-eye/env
WorkingDirectory=/opt/plant-eye
ExecStart=/opt/plant-eye/.venv/bin/uvicorn plant_eye.main:app --host 0.0.0.0 --port 8420
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
